# deploy/systemd/upwork-collector.timer
[Unit]
Description=Run Upwork scheduled collector every 60 minutes

[Timer]
OnActiveSec=60m
OnUnitActiveSec=60m
Persistent=true

[Install]
WantedBy=timers.target
