# deploy/systemd/upwork-collector.timer
[Unit]
Description=Run Upwork scheduled collector every 30 minutes

[Timer]
OnActiveSec=30m
OnUnitActiveSec=30m
Persistent=true

[Install]
WantedBy=timers.target
